# gunicorn/gunicorn.service

Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/envorior
ExecStart=/home/ubuntu/envorior/env/bin/gunicorn \
          --access-logfile - \
          --workers 3 \
          --bind unix:/run/gunicorn.sock \
          mysite.wsgi:application
[Install]
WantedBy=multi-user.target
